# staging.ini
[app:main]
use = egg:RandoPony

pyramid.reload_templates = false
pyramid.debug_authorization = false
pyramid.debug_notfound = false
pyramid.debug_routematch = false
pyramid.default_locale_name = en
pyramid.includes =
    pyramid_deform
    pyramid_mailer
    pyramid_mako
    pyramid_persona
    pyramid_tm

production_deployment = true

# google_drive.username: Set in randopony/private_credentials.py
# google_drive.password: Set in randopony/private_credentials.py

mail.queue_path = %(here)s/mail
mail.host = smtp.webfaction.com
mail.port = 25
# mail.username: Set in randopony/private_credentials.py
# mail.password: Set in randopony/private_credentials.py

mako.directories = randopony:templates

pyramid_deform.template_search_path = randopony:templates/deform

# persona.secret: Set in randopony/private_credentials.py
persona.audiences = http://randopony.bcrandonneur.webfactional.com
persona.siteName = RandoPony

sqlalchemy.url = sqlite:///%(here)s/RandoPony-staging.sqlite

timezone = Canada/Pacific

[server:main]
use = egg:waitress#main
host = 127.0.0.1
port = 23147

[pshell]
m = randopony.models
session = randopony.models.meta.DBSession
t = transaction


# Begin dev environment supervisord configuration

[supervisord]
logfile = %(here)s/supervisord.log

[program:celery]
command = %(here)s/bin/celery worker --config randopony.celery_config --loglevel info
redirect_stderr = true
stdout_logfile = %(here)s/celery.log

# End dev environment supervisord configuration


# Begin logging configuration

[loggers]
keys = root, randopony, sqlalchemy

[handlers]
keys = console

[formatters]
keys = generic

[logger_root]
level = WARN
handlers = console

[logger_randopony]
level = WARN
handlers =
qualname = randopony

[logger_sqlalchemy]
level = WARN
handlers =
qualname = sqlalchemy.engine
# "level = INFO" logs SQL queries.
# "level = DEBUG" logs SQL queries and results.
# "level = WARN" logs neither.  (Recommended for production systems.)

[handler_console]
class = StreamHandler
args = (sys.stderr,)
level = NOTSET
formatter = generic

[formatter_generic]
format = %(asctime)s %(levelname)-5.5s [%(name)s][%(threadName)s] %(message)s

# End logging configuration
